# site-cookbooks/tcpcrypt/files/default/tcpcrypt.conf
#!upstart
description "tcpcrypt"

start on startup
stop on shutdown

respawn

script
   chdir /opt/tcpcrypt/user
   exec /bin/bash launch_tcpcryptd.sh
end script
