# site-cookbooks/tor/templates/default/tor-dns-proxy.upstart.conf.erb
#!upstart
description "proxymachine socks4a resolver for tor"

start on startup
stop on shutdown

respawn

setuid nobody
setgid nogroup

script
  start-stop-daemon --start --make-pidfile --pidfile /var/run/tor-dns-proxy.pid --exec proxymachine -p <%= node["tor"]["socks-resolver-port"] %> -c /opt/tor-dns-proxy.rb
end script
