# site-cookbooks/tor/templates/default/tor-dns-proxy.upstart.conf.erb
#!upstart
description "proxymachine socks4a resolver for tor"

start on startup
stop on shutdown

respawn

script
  start-stop-daemon --start --make-pidfile --pidfile /var/run/tor-dns-proxy.pid --chuid nobody:nogroup --exec /usr/local/bin/proxymachine -- -p <%= node["tor"]["socks-resolver-port"] %> -c /opt/tor-dns-proxy.rb
end script
