# site-cookbooks/tor/templates/default/tor-dns-proxy.upstart.conf.erb
#!upstart
description "proxymachine socks4a resolver for tor"

start on startup
stop on shutdown

respawn

script
  exec proxymachine -p <%= node["tor"]["socks-resolver-port"] %> -c /opt/tor-dns-proxy.rb
end script
