# site-cookbooks/subsonic/templates/default/subsonic.upstart.conf.erb
#!upstart
description "subsonic music server on jetty-runner"

start on startup
stop on shutdown

respawn

script
  start-stop-daemon --start --make-pidfile --pidfile /var/run/subsonic.pid --chuid <%= node["subsonic"]["user"] %>:<%= node["subsonic"]["group"] %> --exec /usr/bin/java -- -jar <%= @jetty_path %> --port <%= node["subsonic"]["port"] %> <%= @subsonic_path %>
end script
