# site-cookbooks/znc/templates/default/socat.conf.erb
#!upstart
description "socat for <%= @network %> irc network over tor"

start on startup
stop on shutdown

respawn

script
  start-stop-daemon --start --make-pidfile --pidfile /var/run/socat-<%= @network %>.pid --chuid nobody:nogroup --exec /usr/bin/socat -- TCP4-LISTEN:<%= @listen_port %>,fork SOCKS4A:localhost:<%= @network_address %>,socksport=<%= node["tor"]["socks-port"] %> > /var/log/socat-<%= @network %>.log
end script
