# site-cookbooks/znc/templates/default/socat-freenode.conf.erb
#!upstart
description "socat for freenode over tor with ssl"

start on startup
stop on shutdown

respawn

setuid nobody
setgid nogroup

script
   socat TCP4-LISTEN:6696,fork SOCKS4A:localhost:$(dig +short irc.tor.freenode.net cname | sed -E "s/([^\.]+\.[^\.]+)\./\1/"):6697,socksport=9050
end script
