# site-cookbooks/znc/templates/default/socat-freenode.conf.erb
#!upstart
description "socat for freenode over tor with ssl"

start on startup
stop on shutdown

respawn

script
  start-stop-daemon --start --make-pidfile --pidfile /var/run/socat-freenode.pid --chuid nobody:nogroup --exec /usr/bin/socat -- TCP4-LISTEN:<%= node["znc"]["socat-freenode-port"] %>,fork SOCKS4A:localhost:$(dig +short irc.tor.freenode.net cname | sed -E "s/([^\.]+\.[^\.]+)\./\1/"):6697,socksport=<%= node["tor"]["socks-port"] %> > /var/log/socat-freenode.log
end script
