# site-cookbooks/dnscrypt-proxy/templates/default/dnscrypt-proxy.conf.erb
#!upstart
description "dnscrypt-proxy-<%= @resolver %>"

start on (local-filesystems and started dbus and stopped udevtrigger)
stop on runlevel [016]

respawn

script
  exec /usr/local/sbin/dnscrypt-proxy --user=<%= node["dnscrypt-proxy"]["user"] %> --local-address=<%= p["local-address"] %> --provider-key=<%= p["provider-key"] %> --provider-name=<%= p["provider-address"] %> --resolver-address=<%= p["resolver-address"] %>:<%= p["resolver-port"] %>
  <% p = node["dnscrypt-proxy"][@resolver] %>
end script
